# supervisor/supervisord.conf
[unix_http_server]
file=/srv/spuelmeister/supervisor/supervisor.sock
chmod=0700

[inet_http_server]
port=0.0.0.0:9001
username=root
password=blafasel

[supervisord]
logfile=/srv/spuelmeister/supervisor/supervisord.log
pidfile=/srv/spuelmeister/supervisor/supervisord.pid
childlogdir=/srv/spuelmeister/supervisor/childlog/
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///srv/spuelmeister/supervisor/supervisor.sock

[eventlistener:spuelmeister]
command=/root/spuelmeister.py
events=PROCESS_STATE,TICK_5

[program:cat]
command=/bin/cat
process_name=%(program_name)s
numprocs=1
directory=/tmp
umask=022
priority=999
autostart=false
autorestart=true

[program:dnsmasq_131]
command=/srv/spuelmeister/supervisor/dnsmasq/dnsmasq-standalone 131
process_name=%(program_name)s
numprocs=1
directory=/srv/spuelmeister/supervisor/dnsmasq
umask=022
priority=999
autostart=false
autorestart=true

[program:dnsmasq_132]
command=/srv/spuelmeister/supervisor/dnsmasq/dnsmasq-standalone 132
process_name=%(program_name)s
numprocs=1
directory=/srv/spuelmeister/supervisor/dnsmasq
umask=022
priority=999
autostart=false
autorestart=true
